# packages/clustermgr4/clustermgr4-4.4.post10.tar.gz/clustermgr4-4.4.post10/clustermgr.service
[Unit]
Description=Gluu Cluster Manager Service

[Service]
Type=forking
ExecStart=/usr/local/bin/clustermgr4-cli start
ExecStop=/usr/local/bin/clustermgr4-cli stop
User=root
Group=root

[Install]
WantedBy=multi-user.target
