# packages/genocide/genocide-84.tar.gz/genocide-84/files/genocide.service
[Unit]
Description=Prosecutor. Court. Reconsider OTP-CR-117/19
After=multi-user.target

[Service]
DynamicUser=True
StateDirectory=genocide
LogsDirectory=genocide
CacheDirectory=genocide
ExecStart=/usr/local/bin/genocided
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
