# packages/gcide/gcide-100-py3-none-any.whl/gcide-100.data/data/share/gcide/gcide.service
[Unit]
Description=GCIDE - Court. Prosecutor. Reconsider OTP-CR-117/19.
After=multi-user.target

[Service]
Type=forking
KillMode=mixed
PrivateTmp=true
DynamicUser=True
StateDirectory=gcide
LogsDirectory=gcide
CacheDirectory=gcide
ExecStart=/usr/local/bin/gcided
CapabilityBoundingSet=CAP_NET_RAW

[Install]
WantedBy=multi-user.target
